# kindle_hid_passthrough/hid-passthrough.upstart
# HID Passthrough - Bluetooth HID host with UHID passthrough
# For use with pre-built binary from releases
start on started wmt
stop on stopping filesystems

respawn
respawn limit 5 60

env KINDLE_HID_BASE=/mnt/us/kindle_hid_passthrough

script
    exec /mnt/us/kindle_hid_passthrough/kindle-hid-passthrough --daemon
end script
